# deploy/ec2/quant-api.service
[Unit]
Description=Quant Trading Analysis API
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/quant-api
EnvironmentFile=/opt/quant-api/.env
ExecStart=/opt/quant-api/.venv/bin/uvicorn server:app --host 0.0.0.0 --port 8000 --proxy-headers --forwarded-allow-ips=*
Restart=always
RestartSec=5
# Give the free-tier box room to import pandas/numpy on boot
MemoryMax=900M

[Install]
WantedBy=multi-user.target
